# user/main.target
[Unit]
Description=main

[Install]
Alias=default.target
